# systemd/clipit.service
[Unit]
Description=Clipit Clipboard manager

[Service]
ExecStart=/usr/bin/clipit --mode 0660

[Install]
WantedBy=multi-user.target
